# data/supervisor/system.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:php-fpm]
command=service php7.1-fpm start
