# data/supervisor/system.conf
[program:nginx]
command=service nginx start

[program:php-fpm]
command=serice php7.1-fpm start
